# litbit.service
[Service]
Type=simple
ExecStart=/usr/local/litbit
User=root
WorkingDirectory=/usr/local/share/litbit/