# litbit-server.service
[Service]
Type=simple
ExecStart=/usr/local/litbit-server
User=root
WorkingDirectory=/usr/local/share/litbit/