# server/litbit-server.service
[Unit]
Description=Run the LitBit web app

[Service]
Type=simple
ExecStart=/usr/local/bin/litbit-server
User=root
WorkingDirectory=/usr/local/share/litbit/

[Install]
WantedBy=multi-user.target
